# user/firstboot-systemctrl/tools/preload.service
[Unit]
Description=Preload execute command

[Service]
ExecStart=/bin/sh /root/firstboot.sh

[Install]
WantedBy=multi-user.target